# Argosss-main/assets/colibri.service
[Unit]
Description=Colibri Daemon — Argos WhisperNode P2P service
Documentation=https://github.com/sigtrip/Argosss
After=network.target

[Service]
Type=simple
User=argos
Group=argos
WorkingDirectory=/opt/argos
ExecStart=/usr/bin/python3 /opt/argos/src/connectivity/colibri_daemon.py \
    --work-dir /var/lib/colibri \
    --pid-file /var/run/colibri.pid
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
